# xob.socket
[Socket]
ListenFIFO=%t/xob.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
