# systemd/my-socks.service
[Unit]
Description=My socks proxy

[Service]
ExecStart=/home/vitja/bin/socks

[Install]
WantedBy=default.target

