# systemd/my-socks.service
[Unit]
Description=My socks proxy
After=network.target

[Service]
ExecStart=/home/vitja/bin/socks

[Install]
WantedBy=default.target
