# systemd/waketv.service
[Unit]
Description=waketv

[Service]
ExecStart=/home/vitja/bin/waketv

[Install]
WantedBy=default.target

