# systemd/xkb-layout.service
[Unit]
Description=XKB Layout Setup

[Service]
ExecStart=/home/vitja/bin/xkb-layout

[Install]
WantedBy=default.target
