# systemd/xkb-layout.service.in
[Unit]
Description=XKB Layout Setup

[Service]
ExecStart=@HOME@/bin/xkb-layout

[Install]
WantedBy=default.target
