# wayland/systemd/sway-nm-applet.service
[Unit]
Description=Sway nm-applet

[Service]
ExecStart=nm-applet --indicator

[Install]
WantedBy=sway-graphical-session.target
