# wayland/systemd/sway-lock.service.in
[Unit]
Description=Sway lock

[Service]
Type=oneshot
ExecStart=@HOME@/bin/sway-lock.sh
