# wayland/systemd/sway-kbd.service
[Unit]
Description=Sway keyboard layout manager

[Service]
ExecStart=swaykbdd

[Install]
WantedBy=sway-graphical-session.target
