# wayland/systemd/sway-wait-sni-ready.service
[Unit]
Description=Wait for SNI to become ready

[Service]
Type=oneshot
ExecStart=/home/vitja/bin/wait-sni-ready
